# script/00-midi-fighter-shutdown.rules
ACTION=="remove", ENV{ID_VENDOR_ID}=="2580", ENV{ID_MODEL_ID}=="0007", RUN+="/sbin/shutdown -h now"
